# .config/systemd/user/hyprland-session.target
# ~/.config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland Session (native, binds graphical-session)
Documentation=man:systemd.special(7)

# Ensure pre target runs first
Wants=graphical-session-pre.target
After=graphical-session-pre.target

# Core binding - allows manual start despite graphical-session.target RefuseManualStart=yes
# BindsTo creates requirement + stop propagation, Before ensures ordering
BindsTo=graphical-session.target
Wants=graphical-session.target
Before=graphical-session.target

# Modern desktop integration - pulls autostart and portals on Arch bleeding-edge
Wants=xdg-desktop-autostart.target
Wants=xdg-document-portal.service
Wants=xdg-permission-store.service
After=xdg-permission-store.service

[Install]
# Intentionally no WantedBy - started manually from autostart.lua:
# systemctl --user start hyprland-session.target
# Services that should run inside Hyprland: PartOf=graphical-session.target
# and WantedBy=graphical-session.target or hyprland-session.target
