# user_scripts/tts_stt/dusky_parakeet/dusky-stt.service
[Unit]
Description=Dusky STT v6 - Unified #1 (5.91% WER) Realtime typing, D3-cold safe
After=pipewire.service pipewire-pulse.service graphical-session.target xdg-desktop-portal.service
Wants=pipewire.service
PartOf=graphical-session.target
StartLimitBurst=5
StartLimitIntervalSec=90

[Service]
Type=simple
ExecStart=%h/contained_apps/uv/dusky_stt_v2/.venv/bin/python %h/contained_apps/uv/dusky_stt_v2/dusky_main.py --daemon
WorkingDirectory=%h/contained_apps/uv/dusky_stt_v2
Environment=HF_HUB_CACHE=%h/.cache/huggingface
Environment=HF_HUB_ENABLE_HF_TRANSFER=1
Environment=PYTHONUNBUFFERED=1
Environment=CUDA_MODULE_LOADING=LAZY
Environment=LD_LIBRARY_PATH=/opt/cuda/lib64
MemoryHigh=32G
MemoryMax=48G
MemorySwapMax=8G
Restart=on-failure
RestartSec=2
TimeoutStopSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dusky-stt

[Install]
WantedBy=default.target
